# scripts/com.reddit.analyzer.historical.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.reddit.analyzer.historical</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jaredteerlink/repos/reddit-analyzer/scripts/cron_historical_collection.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Run every 2 days at 2:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>0</integer> <!-- Sunday -->
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer> <!-- Tuesday -->
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer> <!-- Thursday -->
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>6</integer> <!-- Saturday -->
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/jaredteerlink/repos/reddit-analyzer/logs/launchagent_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jaredteerlink/repos/reddit-analyzer/logs/launchagent_stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/jaredteerlink/repos/reddit-analyzer</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
